# bin/initscript-voteatx-svc
#!/bin/sh
### BEGIN INIT INFO
# Provides:          voteatx-svc
# Required-Start:    $network $syslog
# Required-Stop:     $network $syslog
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Vote ATX service
# Description:       Start the Ruby RACK server that runs the Vote ATX web service.
### END INIT INFO

### Begin config items
#
HOME="/home/ubuntu"
APP_ROOT="$HOME/voteatx-svc"
APP_USER="ubuntu:ubuntu"
SVC="voteatx-svc"
PROG="$APP_ROOT/bin/voteatx-svc"
#
### End config items

USAGE="usage: $0 {start|stop|restart|status|start}"

die() {
	echo "$@" >&2
	exit 1
}

do_start() {
	do_status >/dev/null 2>&1
	if [ $? -eq 0 ] ; then
		die "$0: Service $SVC already running"
	fi

	echo "Starting $SVC ... \c"
	daemon -v \
		--env="APP_ROOT=$APP_ROOT" \
		--env="RACK_ENV=production" \
		--env="HOME=$HOME" \
		--env="PATH=$PATH" \
		--chdir="$APP_ROOT" \
                --user="$APP_USER" \
		--pidfile="$APP_ROOT/log/$SVC.pid" \
		--output="$APP_ROOT/log/$SVC.run" \
		--name="$SVC" \
		bundle exec $PROG 
	echo "OK"
}

do_stop() {
	echo "Stoping $SVC ... \c"
	daemon -v \
		--pidfile="$APP_ROOT/log/$SVC.pid" \
		--name="$SVC" \
 		--stop
	echo "OK"
}

do_status() {
	 daemon -v \
		--pidfile="$APP_ROOT/log/$SVC.pid" \
		--name="$SVC" \
 		--running
}

if [ $# -ne 1 ] ; then
	die "$USAGE"
fi

case "$1" in
	start)	do_start ;;
	stop)	do_stop ;;
	restart) do_stop ; do_start ;;
	status)	do_status ;;
	*)	die "$USAGE" ;; 
esac

